# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create wxlistener user if it doesn't exist
        if ! getent passwd wxlistener >/dev/null; then
            adduser --system --quiet --home /var/lib/wxlistener --no-create-home \
                --shell /bin/false --group --gecos "wxlistener Weather Station Monitor" wxlistener
        fi

        # Create directories
        install -d -m 755 -o root -g root /opt/wxlistener
        install -d -m 755 -o wxlistener -g wxlistener /etc/wxlistener
        install -d -m 755 -o wxlistener -g wxlistener /var/lib/wxlistener

        # Set permissions on config file if it exists
        if [ -f /etc/wxlistener/wxlistener.toml ]; then
            chown root:wxlistener /etc/wxlistener/wxlistener.toml
            chmod 640 /etc/wxlistener/wxlistener.toml
        fi

        # Ensure binary is executable
        if [ -f /opt/wxlistener/wxlistener ]; then
            chown root:root /opt/wxlistener/wxlistener
            chmod 755 /opt/wxlistener/wxlistener
        fi

        # Post-install instructions
        echo ""
        echo "=================================================="
        echo "  wxlistener installed successfully!"
        echo "=================================================="
        echo ""
        echo "Next steps:"
        echo "  1. Edit the configuration file:"
        echo "     sudo nano /etc/wxlistener/wxlistener.toml"
        echo ""
        echo "  2. Enable the service to start on boot:"
        echo "     sudo systemctl enable wxlistener"
        echo ""
        echo "  3. Start the service:"
        echo "     sudo systemctl start wxlistener"
        echo ""
        echo "  4. Check the service status:"
        echo "     sudo systemctl status wxlistener"
        echo ""
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
